# mkosi.images/llamacpp/mkosi.extra/usr/lib/systemd/system/llama-cpp.socket
[Unit]
Description=Socket for llama.cpp server

[Socket]
ListenStream=8000
Service=llama-cpp-proxy.service

[Install]
WantedBy=sockets.target
