# mkosi.extra.bootable/usr/lib/systemd/system/systemd-repart.service.d/no-discard.conf
[Service]
Environment="SYSTEMD_REPART_MKFS_OPTIONS_BTRFS=-K"
