# mkosi.extra.bootable/usr/share/factory/etc/profile.d/direnv.sh
# Enable direnv
eval "$(direnv hook bash)"
