# mkosi.images/atomic/mkosi.conf.d/lenovo-x13s/mkosi.extra/usr/lib/systemd/system/hci-btaddress@.service
[Unit]
Description=HCI bluetooth address fix

[Service]
Type=simple
ExecStart=/bin/sh -c 'sleep 5 && yes | btmgmt -i %I public-addr F4:A8:0D:30:9D:8B'
